# proxy/keyclaw.service
[Unit]
Description=KeyClaw — Local MITM Proxy that strips secrets from LLM traffic
After=network.target
Wants=network.target

[Service]
Type=simple
User=proxy
Restart=always
RestartSec=5
ExecStart=/home/proxy/.cargo/bin/keyclaw proxy start --foreground
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
